# etc/systemd/system/sockd.service.d/after-wg0.conf
[Unit]
After = systemd-networkd-wait-online@wg0.service
Requires = systemd-networkd-wait-online@wg0.service
